# profile.d/aws.sh
# https://github.com/serverless/serverless/issues/3833#issuecomment-398642431
# sls deploy "Profile blah does not exist" problem
export AWS_SDK_LOAD_CONFIG=1
