# profile.d/haskell.sh
# Haskell stuff
append_path_if_exists ${HOME}/.cabal/bin

