# profile.d/linuxbrew.sh
prepend_path_if_exists            \
  /home/linuxbrew/.linuxbrew/bin  \
  /home/linuxbrew/.linuxbrew/sbin

