# profile.d/local_bin.sh
# TODO not exacty sure what chooses to install into .local/bin
#
append_path_if_exists ~/.local/bin
