# profile.d/vagrant.sh
# TODO this is linux only ?
export VAGRANT_DEFAULT_PROVIDER=libvirt
